# files/portfwd-natpmp@22.timer
[Unit]
Description=NAT-PMP port fwd interval

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
