# files/portfwd-upnp@22.timer
[Unit]
Description=run upnp port fwd interval

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
