# files/sunday-motion.service
# /etc/systemd/system/sunday-motion.service
# Automatically generated by systemd-sysv-generator
# but then modified by john to be like the github repo suggestoin

[Unit]
Description=sundaycam motion daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/motion -n
User=motion
PIDFile=/var/run/motion.pid
RestartSec=5

[Install]
WantedBy=multi-user.target
